# .config/autostart/ibus-daemon.desktop
[Desktop entry]
Type=Apllication
Name=IBus Daemon
Exec=ibus-daemon -d
